# mvmm.service
[Unit]
Description=Multi-VM Monitor (mvmm)
After=network.target

[Service]
ExecStart=/usr/local/bin/mvmm
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target